# lib/systemd/system/revpi-hciuart@.service
[Unit]
Description=Revolution Pi Bluetooth Device %i
After=systemd-udev-settle.service

[Service]
Type=forking
ExecStart=/usr/bin/revpi-btuart %i
